# crontab
*/2 * * * * python /app/create_real_article.py >> /var/log/cron.log 2>&1
